# license-verificator.service
[Unit]
Description=License-Verificator
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/programs/License-Verificator
Environment=COMPOSE_HTTP_TIMEOUT=3600
ExecStart=/usr/bin/env /usr/bin/docker-compose -f /home/ubuntu/programs/License-Verificator/docker-compose.yaml up --build
ExecStop=/usr/bin/env /usr/bin/docker-compose -f /home/ubuntu/programs/License-Verificator/docker-compose.yaml stop
StandardOutput=syslog
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

